# config/systemd/user/obsidian-bisync.service
[Unit]
Description=Obsidian Bi-directional Sync

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gdrive-bisync %h/Documents/obsidian gdrive:obsidian
